# env/chef/cookbooks/pluto/templates/mac_os_x/pluto_init.plist.erb
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>cc.mrtesla.pluto-init</string>
  <key>ServiceDescription</key>
  <string>pluto-init - starts and monitors a collection of runsv(8) processes</string>

  <key>QueueDirectories</key>
  <array>
    <string><%= @alice_root %>/var/runit/enabled-services</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin</string>
  </dict>

  <key>OnDemand</key>
  <false/>

  <key>Program</key>
  <string><%= @alice_root %>/bin/pluto-init</string>
</dict>
</plist>
